# docker/modern/rootfs/etc/s6-overlay/s6-rc.d/init-runtime/up
#!/command/execlineb -P
foreground { mkdir -p /tmp/runtime }
foreground { chmod 0700 /tmp/runtime }
