# docker/modern/rootfs/etc/s6-overlay/s6-rc.d/pipewire/run
#!/command/execlineb -P
with-contenv
export XDG_RUNTIME_DIR /tmp/runtime
exec pipewire
