# docker/modern/rootfs/etc/s6-overlay/s6-rc.d/game/run
#!/command/execlineb -P
with-contenv
/bin/bash -c "
exec 2>&1

# Wait for X server to be ready
echo 'Waiting for X server...'
while ! xdpyinfo -display :99 >/dev/null 2>&1; do
    sleep 0.5
done
echo 'X server ready'

exec /opt/bin/start-game.sh
"
